# docker-vhpc/rocky10/docker/head/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
stdout_logfile=/var/log/sshd.out
stderr_logfile=/var/log/sshd.err
