# reference/slurm/storage/supervisord.conf
[supervisord]
nodaemon=true

[program:storage-start]
command=/usr/local/bin/start-storage.sh
autostart=true
stderr_logfile=/var/log/storage.err
stdout_logfile=/var/log/storage.out
